# webhook.service
[Unit]
Description=GitHub Webhook Server for Gatsby Blog
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/gatsby-blog
ExecStart=/usr/bin/node /var/www/gatsby-blog/webhook.js
Restart=on-failure
Environment=WEBHOOK_PORT=9000
Environment=WEBHOOK_BRANCH=main

[Install]
WantedBy=multi-user.target 